# deployment/docker/standalone/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/bin/bash gunicorn_starter.bash
